# etc/cron.d/rig
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# m h dom mon dow user	command
#@reboot	root	/home/user/git/linux_mining/x11vnc.sh
@reboot		root	sleep 15 && tmux new-session -d -s miner '/bin/bash /home/user/miner/my_start_miner.bash'
@reboot		root	sleep 15 && /usr/bin/python3 /home/user/miner/overclock.py
#
