# user/.config/autostart/x11vnx.desktop
[Desktop Entry]
Type=Application
Exec=/home/user/x11vnc.sh
Hidden=false
X-GNOME-Autostart-enabled=true
Name=x11vnc
Comment=x11vnc server